# sources/sacaa/deploy/sacaa-cycle.service
# /etc/systemd/system/sacaa-cycle.service
[Unit]
Description=SACAA AIID (South Africa) ingest cycle (weekly; discover->fetch->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/sacaa
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/sacaa/sacaa-cycle.log
StandardError=append:/opt/sacaa/sacaa-cycle.log
TimeoutStartSec=14400
