# sources/otkes/deploy/otkes-cycle.timer
# /etc/systemd/system/otkes-cycle.timer
[Unit]
Description=Weekly OTKES (Finland) ingest

[Timer]
OnCalendar=Sun 15:30 UTC
Persistent=true

[Install]
WantedBy=timers.target
