# sources/otkes/deploy/otkes-cycle.service
# /etc/systemd/system/otkes-cycle.service
[Unit]
Description=OTKES (Finland) ingest cycle (weekly; Playwright/Xvfb; discover->fetch->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/otkes
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/otkes/otkes-cycle.log
StandardError=append:/opt/otkes/otkes-cycle.log
TimeoutStartSec=21600
