# sources/shk/deploy/shk-cycle.timer
# /etc/systemd/system/shk-cycle.timer
[Unit]
Description=Weekly SHK (Sweden) ingest cycle

[Timer]
OnCalendar=Sun 13:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
