# sources/gpiaaf/deploy/gpiaaf-cycle.timer
# /etc/systemd/system/gpiaaf-cycle.timer
[Unit]
Description=Weekly GPIAAF (Portugal) civil-aviation ingest

[Timer]
OnCalendar=Sun 14:30 UTC
Persistent=true

[Install]
WantedBy=timers.target
