# sources/taic/deploy/taic-cycle.timer
# /etc/systemd/system/taic-cycle.timer
[Unit]
Description=Weekly TAIC (New Zealand) ingest cycle

[Timer]
OnCalendar=Sun 19:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
