# sources/ciaado/deploy/ciaado-cycle.service
# /etc/systemd/system/ciaado-cycle.service
[Unit]
Description=CIAA (Dominican Republic) ingest cycle (weekly; discover->fetch->parse->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/ciaado
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/ciaado/ciaado-cycle.log
StandardError=append:/opt/ciaado/ciaado-cycle.log
TimeoutStartSec=10800
