# sources/griaa/deploy/griaa-cycle.service
# /etc/systemd/system/griaa-cycle.service
[Unit]
Description=GRIAA (Colombia) ingest cycle (weekly; discover->fetch->parse->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/griaa
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/griaa/griaa-cycle.log
StandardError=append:/opt/griaa/griaa-cycle.log
TimeoutStartSec=10800
