# sources/tsb/deploy/tsb-cycle.service
# /etc/systemd/system/tsb-cycle.service
[Unit]
Description=TSB (Canada) ingest cycle (weekly; discover->fetch->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/tsb
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/tsb/tsb-cycle.log
StandardError=append:/opt/tsb/tsb-cycle.log
TimeoutStartSec=7200
