# sources/aibdk/deploy/aibdk-cycle.timer
# /etc/systemd/system/aibdk-cycle.timer
[Unit]
Description=Weekly AIB Denmark ingest cycle

[Timer]
OnCalendar=Sun 07:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
