# sources/knkt/deploy/knkt-cycle.service
# /etc/systemd/system/knkt-cycle.service
[Unit]
Description=KNKT NTSC (Indonesia) ingest cycle (weekly; discover->fetch->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/knkt
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/knkt/knkt-cycle.log
StandardError=append:/opt/knkt/knkt-cycle.log
TimeoutStartSec=7200
