# sources/knkt/deploy/knkt-cycle.timer
# /etc/systemd/system/knkt-cycle.timer
[Unit]
Description=Weekly KNKT NTSC (Indonesia) ingest cycle

[Timer]
OnCalendar=Sun 23:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
