# sources/india/deploy/india-cycle.timer
# /etc/systemd/system/india-cycle.timer
[Unit]
Description=Weekly AAIB India ingest cycle

[Timer]
OnCalendar=Sun 21:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
