# sources/ansv/deploy/ansv-cycle.service
# /etc/systemd/system/ansv-cycle.service
[Unit]
Description=ANSV (Italy) ingest cycle (weekly; discover->fetch->parse->build)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/ansv
ExecStart=/bin/bash -lc "deploy/run-cycle.sh"
StandardOutput=append:/opt/ansv/ansv-cycle.log
StandardError=append:/opt/ansv/ansv-cycle.log
TimeoutStartSec=10800
